# crontab
*/1 * * * * root curl "http://localhost:5000/test" >> /var/log/aarc_cron.log 2>&1 
# Don't remove the empty line at the end of this file. It is required to run the cron job
